# systemd/fanuc-monitor.service
[Unit]
Description=FANUC CNC Monitor Service
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/iskra-fanuc-hub
Environment="PATH=/home/admin/iskra-fanuc-hub/venv/bin"
ExecStart=/home/admin/iskra-fanuc-hub/venv/bin/python -m app.main /home/admin/iskra-fanuc-hub/config.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
